# dovecot/root/etc/cron.d/dovecot
* * * * * root curl "http://solr:8983/solr/dovecot/update?commit=true" &>/dev/null
