# spamd/root/etc/supervisor/conf.d/spamd.conf
[program:spamd]
command=/usr/local/sbin/spamd.sh
autorestart=true

[program:spamass-milter]
command=/usr/sbin/spamass-milter -p inet:1234 -m
user=spamass-milter
autorestart=true
