# rspamd/root/etc/supervisor/conf.d/rspamd.conf
[program:rspamd]
command=/usr/sbin/rspamd -u rspamd -g rspamd -p /var/run/rspamd.pid -f
autorestart=true
