# dovecot/root/etc/supervisor/conf.d/dovecot.conf
[program:dovecot]
command=/usr/sbin/dovecot -F -c /etc/dovecot/dovecot.conf
autorestart=true
