# dovecot/root/etc/periodic/15min/dovecot
#!/bin/sh
curl 'http://solr:8983/solr/dovecot/update?commit=true' >/dev/null 2>&1
