# opendkim/root/etc/supervisor/conf.d/opendkim.conf
[program:opendkim]
command=/usr/sbin/opendkim -f
autorestart=true
stdout_logfile=syslog
stderr_logfile=syslog
